# scripts/deploy/templates/linux/riskhub-scheduler.service.tmpl
[Unit]
Description=RiskHub Backend Scheduler
After=network-online.target __REDIS_SERVICE__.service
Wants=network-online.target
Requires=__REDIS_SERVICE__.service

[Service]
Type=simple
User=riskhub
Group=riskhub
WorkingDirectory=__CURRENT_LINK__/backend
EnvironmentFile=__RUNTIME_ENV__
Environment=ENABLE_SCHEDULER=true
ExecStart=__CURRENT_LINK__/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
